# web3/config/history-archive.conf
; override to disable history-archive
[program:history-archive]
autostart=false
command=/bin/true